# conf/supervise-nginx.conf
[program:Supervise_nginx]
user=http
command=/usr/sbin/nginx -f /etc/nginx/nginx.conf
redirect_stderr=true
autorestart=true
