# conf/supervise-cron.conf
[program:Supervise_cron]
user=cron
command=/usr/sbin/cron
autorestart=true
